# ansible/roles/amneziawg/files/awg-hub-router.service
[Unit]
Description=Restore per-user AWG hub routing on boot
After=awg-quick@awg0.service awg-quick@awg1.service awg-quick@awg2.service
Wants=awg-quick@awg0.service awg-quick@awg1.service awg-quick@awg2.service

[Service]
Type=oneshot
ExecStart=/etc/amnezia/amneziawg/switch-hub.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
